# MovableTypeRebuilder/scheduler/com.user.movabletype-rebuilder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.movabletype-rebuilder</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder && source venv/bin/activate && python scripts/mt_rebuilder.py</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>1</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder/logs/launchd_stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>MT_SITE_URL</key>
        <string>https://your-movabletype-site.com</string>
        <key>MT_USERNAME</key>
        <string>your-username</string>
        <key>MT_PASSWORD</key>
        <string>your-password</string>
        <key>MT_BLOG_ID</key>
        <string>1</string>
        <key>MT_SITE_NAME</key>
        <string>MovableTypeサイト</string>
        <key>EMAIL_USERNAME</key>
        <string>your-email@gmail.com</string>
        <key>EMAIL_PASSWORD</key>
        <string>your-app-password</string>
        <key>FROM_EMAIL</key>
        <string>your-email@gmail.com</string>
        <key>TO_EMAIL</key>
        <string>your-email@gmail.com</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
