# MovableTypeRebuilder/scheduler/com.user.movabletype-rebuilder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.movabletype-rebuilder</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder/venv/bin/python</string>
        <string>scripts/mt_rebuilder.py</string>
        <string>--schedule</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>1</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder/logs/launchd_stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/takuhito/NotionWorkflowTools/MovableTypeRebuilder</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
